# deploy/crontab
# most reports
10 5 * * * /home/ubuntu/daily.sh > /home/ubuntu/daily.log 2>&1

# realtime reports
*/1 * * * * /home/ubuntu/realtime.sh > /home/ubuntu/realtime.log 2>&1
